# bin/drive-mirroring.service
[Unit]
Description=Drive sync

[Service]
Type=simple
ExecStart=/home/frain/bin/drive-mirroring
