# bin/drive-mirroring.timer
[Unit]
Description=Drive sync

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target